# configurations/supervisor.d/entrypoint.ini
[program:entrypoint]
user = root
priority = 0
command = bash -c "/usr/local/bin/entrypoint.sh"
redirect_stderr = true
