# configurations/supervisor.d/sshd.ini
[program:sshd]
user = root
autostart = true
command = /usr/sbin/sshd -D -E /var/log/sshd.log -f /etc/ssh/sshd_config
stdout_logfile = /var/log/sshd-common.log
stderr_logfile = /var/log/sshd-out.log
