# configurations/supervisor.d/sshd.ini
[program:sshd]
user = root
autostart = true
priority = 1
command = /usr/sbin/sshd -D -E /var/log/sshd.log -f /etc/ssh/sshd_config
redirect_stderr = true
