# configurations/supervisor.d/wildfly.ini
[program:wildfly]
user = wildfly
autostart = true
command = /opt/wildfly/bin/standalone.sh -b 0.0.0.0 -bmanagement=0.0.0.0
stdout_logfile = /var/log/wildfly-common.log
stderr_logfile = /var/log/wildfly-error.log
